# intelligence/com.cashclaw.alphaear.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.cashclaw.alphaear</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>server.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/you/algo-trader/intelligence</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>LLM_FAST_TRIAGE_URL</key>
        <string>http://127.0.0.1:11436/v1</string>
        <key>LLM_FAST_TRIAGE_MODEL</key>
        <string>mlx-community/NVIDIA-Nemotron-3-Nano-30B-A3B-4bit</string>
        <key>SENTIMENT_MODE</key>
        <string>bert</string>
        <key>SIDECAR_PORT</key>
        <string>8100</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/alphaear-sidecar.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/alphaear-sidecar.err</string>
</dict>
</plist>
